# debian-bookworm/.devcontainer/config/ollama.conf
[program:ollama]
command=/usr/local/bin/ollama serve
user=ollama
directory=/var/lib/ollama
autostart=true
autorestart=true
stderr_logfile=/var/log/ollama.err
stdout_logfile=/var/log/ollama.log
environment=HOME="/var/lib/ollama",USER="ollama",OLLAMA_MODELS="/var/lib/ollama/models"